# overlay/etc/xdg/autostart/octolinux-session.desktop
[Desktop Entry]
Type=Application
Name=OctoLinux Session Setup
Exec=/usr/local/bin/octolinux-first-login
NoDisplay=true
